# soc.ini
[program:clawer]
autorestart=True
autostart=True
redirect_stderr=True
command=python clawer.py
user=betta
directory=/home/betta/sharp
stdout_logfile=/tmp/sharp_clawer_stdout.log

[program:list]
autorestart=True
autostart=True
redirect_stderr=True
command=python list.py
user=betta
directory=/home/betta/sharp
stdout_logfile=/tmp/sharp_list_stdout.log

[program:pipeline]
autorestart=True
autostart=True
redirect_stderr=True
command=python pipeline.py
user=betta
directory=/home/betta/sharp
stdout_logfile=/tmp/sharp_pipeline_stdout.log

[program:result]
autorestart=True
autostart=True
redirect_stderr=True
command=python result.py
user=betta
directory=/home/betta/sharp
stdout_logfile=/tmp/sharp_result_stdout.log
